# packages/video-worker/root/etc/s6-overlay/s6-rc.d/svc-videoworker/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "Starting video worker.."

cd /home/node/video-worker
npx ts-node --transpile-only src/index.ts
